# docker/8.1/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:php]
command=%(ENV_SUPERVISOR_PHP_COMMAND)s
user=%(ENV_SUPERVISOR_PHP_USER)s
environment=LARAVEL_SAIL="1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:scheduler]
command=/bin/sh -c "while [ true ]; do (/usr/bin/php -d variables_order=EGPCS /var/www/html/artisan schedule:run --verbose --no-interaction &); sleep 60; done"
user=%(ENV_SUPERVISOR_PHP_USER)s
environment=LARAVEL_SAIL="1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:worker-certification]
command=/usr/bin/php -d variables_order=EGPCS /var/www/html/artisan queue:work --queue=CreateCertificationUserJob,InactivateCertificateQualifiedJob --timeout=99999
user=%(ENV_SUPERVISOR_PHP_USER)s
environment=LARAVEL_SAIL="1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:default]
command=/usr/bin/php -d variables_order=EGPCS /var/www/html/artisan queue:work --queue=default,auth.events
user=%(ENV_SUPERVISOR_PHP_USER)s
environment=LARAVEL_SAIL="1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:email]
command=/usr/bin/php -d variables_order=EGPCS /var/www/html/artisan queue:work --queue=email,mail,notification
user=%(ENV_SUPERVISOR_PHP_USER)s
environment=LARAVEL_SAIL="1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:events]
command=/usr/bin/php -d variables_order=EGPCS /var/www/html/artisan queue:work --queue=events.user.added
user=%(ENV_SUPERVISOR_PHP_USER)s
environment=LARAVEL_SAIL="1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:ldap]
command=/usr/bin/php -d variables_order=EGPCS /var/www/html/artisan queue:work --queue=CreateLdapUserJob,LockLdapUserJob,UnlockLdapUserJob,ResetLdapPasswordJob,ldap.events,ldap.jobs,events.ldap.bind,UpdateManagerAndGroupJob,UpdateTelephoneJob,NotifyLinkedCollaboratorsToDisabledManagersJob
user=%(ENV_SUPERVISOR_PHP_USER)s
environment=LARAVEL_SAIL="1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:outsourced]
command=/usr/bin/php -d variables_order=EGPCS /var/www/html/artisan queue:work --queue=VerifyOutsourcedValidationTimeJob,DisableOutsourcedValidationTimeExpiredJob,ValidateOutsourcedNotRegisteredJob,InactivateOutsourcedNotRegisteredJob --timeout=99999
user=%(ENV_SUPERVISOR_PHP_USER)s
environment=LARAVEL_SAIL="1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:soulmv]
command=/usr/bin/php -d variables_order=EGPCS /var/www/html/artisan queue:work --queue=CreateSoulMVUserJob,LockSoulMVUserJob,UnlockSoulMVUserJob,NotifyManagerPasswordChangedJob,NotifyResetPasswordSoulMVJob,soulmv-events,soulmv.jobs
user=%(ENV_SUPERVISOR_PHP_USER)s
environment=LARAVEL_SAIL="1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:worker-sync-jobs]
command=/usr/bin/php -d variables_order=EGPCS /var/www/html/artisan queue:work --queue=SyncAllCollaboratorsToUsersJob,SyncAllUsersToRolesJob,SyncColaboratorsConditionToUsersJob,SyncAllCollaboratorsRegistryToLdapJob,SyncAllCollaboratorsRolesToLdapJob,RemoveAllLdapAndSoulMVUserPermissionsJob --timeout=99999
user=%(ENV_SUPERVISOR_PHP_USER)s
environment=LARAVEL_SAIL="1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:worker-user]
command=/usr/bin/php -d variables_order=EGPCS /var/www/html/artisan queue:work --queue=UpdateCollaboratorUserJob,UpdateOutsourcedUserJob,UpdateOutsourcedCategoryJob
user=%(ENV_SUPERVISOR_PHP_USER)s
environment=LARAVEL_SAIL="1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
